# extra-repos/systemd/postmarketos-base-systemd/rootfs-usr-lib-systemd-system-apk-inhibit.service
[Unit]
Description=Hold sleep inhibitor while apk is running

[Service]
Type=simple
ExecStart=/usr/libexec/apk-inhibit
