# extra-repos/systemd/postmarketos-base-systemd/rootfs-usr-lib-systemd-system-apk-inhibit.path
[Unit]
Description=Watch for apk lock file

[Path]
PathExists=/usr/lib/apk/db/lock
Unit=apk-inhibit.service

[Install]
WantedBy=multi-user.target
